# debian/uif.postrm
#!/bin/sh

set -e

if [ "$1" = purge ]; then
	update-rc.d -f uif remove > /dev/null
fi

#DEBHELPER#

exit 0
