# debian/uif.prerm
#!/bin/sh

set -e

# only stop firewall when uif is going to be removed
if [ "$1" = "remove" -o "$1" = "deconfigure" ]; then
	invoke-rc.d uif stop
else
	echo "Please restart uif manually using: invoke-rc.d uif restart"
fi

#DEBHELPER#

exit 0
